# uxsite/deploy/supervisord.conf
; /etc/supervisor/conf.d/uxsite.conf
; supervisor program definition for the uxsite Hotframe app.
; Reload with:  sudo supervisorctl reread && sudo supervisorctl update
; Restart  with: sudo supervisorctl restart uxsite

[program:uxsite]
; The Lightsail Nano plan ships only 512 MB RAM; Hotframe imports
; SQLAlchemy + Pydantic + datastar + uvicorn workers, which adds up
; quickly. One worker keeps the steady-state RSS around ~150 MB.
; Upgrade to micro (1 GB) to comfortably run 2+ workers.
command=/srv/uxsite/.venv/bin/gunicorn uxsite.main:app
    --bind 127.0.0.1:8000
    --workers 1
    --worker-class uvicorn.workers.UvicornWorker
    --timeout 60
    --access-logfile -
    --error-logfile -
directory=/srv/uxsite
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=15
environment=PORT="8000",PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/uxsite/stdout.log
stderr_logfile=/var/log/uxsite/stderr.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
